# displaypi.service
[Unit]
Description=DisplayPi kiosk browser
After=network.target

[Service]
Type=simple
User=pi
Environment=DISPLAY=:0
Environment=XDG_RUNTIME_DIR=/run/user/1000
ExecStart=/usr/bin/python3 /home/pi/DisplayPi/displaypi.py
Restart=on-failure

[Install]
WantedBy=graphical.target
